# maxwell_defs/upstart/maxwell.conf
description	"bringup maxwell"

start on net-device-up IFACE=wlan0
stop on net-device-down IFACE=wlan0

console output
respawn

exec maxwell-start
pre-stop exec maxwell-stop
